# server/systemd/minecraft-collect-node-snapshot.service
[Unit]
Description=Collect node snapshot (proc + RCON list + spark health)
After=minecraft.service

[Service]
Type=oneshot
ExecStart=/opt/minecraft/bin/collect-node-snapshot.sh
