# server/systemd/minecraft.service
[Unit]
Description=Minecraft (Fabric) server
# ネットワークが整ってから起動する（EIP関連付け/S3同期はinstall.shが先行実施済み）。
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/minecraft.env
WorkingDirectory=/opt/minecraft/server
# Fabric Server Launcher jar は固定名 fabric-server-launch.jar とする（jar名規約。docs参照）。
# -Xmx6G/-Xms6G は c8g.xlarge(8GB) 前提。
# JDK 25 向けチューニング（旧 Aikar G1 系フラグから刷新）:
#   - UseZGC: Generational ZGC は Java 23+ で ZGC のデフォルト動作、Java 25 で完全安定。
#     チャンク生成等の長停止 (G1 200ms 目標) を <10ms に抑える。
#   - UseCompactObjectHeaders (JEP 519, Java 25 stable): オブジェクトヘッダを 12B → 8B に圧縮。
#     BlockPos/ItemStack/Entity 等の大量小オブジェクトで実効ヒープを 10-15% 節約しキャッシュ効率も改善。
#   - UseStringDeduplication: ZGC と組み合わせて重複 String の内部 byte[] を共有し追加でヒープ節約
#     （ブロック ID / NBT キー等が大量に重複するため Minecraft では効きやすい）。
#   - DisableExplicitGC: MOD 由来の System.gc() を無視し ZGC のリズムを乱さない。
#   - PerfDisableSharedMem: /tmp/hsperfdata の競合を回避。
#   - -XX:StartFlightRecording: JVM 公式の構造化イベント録画 (JFR) を継続稼働させ、
#     後段の collect-jfr.sh が 60s ごとに jcmd JFR.dump → jfr print --json → JSONL 化して
#     S3+Athena に流す。settings=default で overhead ~1%。disk=true で disk-backed の
#     continuous recording、maxage=10m / maxsize=100M で過去 10 分または 100MB に
#     収まる範囲を NVMe ephemeral 上にロールし続ける。
#     旧 -Xlog:gc* + parse-gc-log.sh (正規表現方式) はフォーマット揺れとイベント
#     取りこぼしで品質が悪く JDK 版互換性も保証できないため JFR に統一した。
# JDK 25 で自動適用される改善（追加フラグ不要）: C2 JIT 最適化, JFR cooperative sampling 等。
#
# 撤去フラグ:
#   - AlwaysPreTouch / UseTransparentHugePages の組み合わせは、6GB ヒープに対して
#     Linux の khugepaged (THP 整列スレッド) が defrag に走り 1 コアを占有 → ページ操作で
#     待たされた全プロセス (Java / SSM Agent / CWAgent) が同時 hang する事象を実環境で
#     観測したため撤去。インスタンスは ping 通るがすべての user-space が応答停止する。
#     ZGC は内部で必要に応じて大ページを扱う設計なので、両フラグなしでも性能影響は微小。
ExecStart=/usr/bin/java \
    -Xmx6G -Xms6G \
    -XX:+UseZGC \
    -XX:+UseCompactObjectHeaders \
    -XX:+UseStringDeduplication \
    -XX:+DisableExplicitGC \
    -XX:+PerfDisableSharedMem \
    -XX:StartFlightRecording=name=mc,filename=/opt/minecraft/server/metrics/mc.jfr,disk=true,maxage=10m,maxsize=100M,settings=default \
    -jar fabric-server-launch.jar nogui
# 停止時（OSシャットダウン/ASG終了/スポット中断いずれもACPI停止→systemd停止）に必ずsaveする。
# 続けて metrics の最終 flush を実行する (`-` prefix: 観測 flush が失敗しても停止フローは止めない)。
ExecStop=/opt/minecraft/bin/save-and-sync.sh
ExecStop=-/opt/minecraft/bin/flush-metrics.sh
# スポット中断猶予(約2分)内に save+S3同期+metrics flush を終えるための停止タイムアウト。
TimeoutStopSec=110
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
