# server/systemd/minecraft-collect-jfr.service
[Unit]
Description=Collect JFR snapshot and emit GC events JSONL
After=minecraft.service

[Service]
Type=oneshot
ExecStart=/opt/minecraft/bin/collect-jfr.sh
