# server/systemd/world-sync.timer
[Unit]
Description=Periodic Minecraft world sync to S3 (every 15 min)

[Timer]
# 起動15分後に初回、以降は前回実行から15分間隔で定期同期する。
OnUnitActiveSec=15min
OnBootSec=15min
Unit=world-sync.service

[Install]
WantedBy=timers.target
